# doc/Start Odoo Automatically On Boot

sudo cp /opt/odoo/ctlscript.sh /etc/init.d/odoo



sudo pico  /etc/init.d/odoo

Add  the following lines at the beginning of the /etc/init.d/odoo

#!/bin/sh
### BEGIN INIT INFO
# Provides:          odoo
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog  
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: odoo
# Description:       odoo
### END INIT INFO

 


sudo chmod +x /etc/init.d/odoo

sudo update-rc.d -f odoo defaults
sudo update-rc.d -f odoo enable


sudo service odoo start


##############################################################
To remove the banner link to the Bitnami Info page, follow these steps:

sudo /opt/odoo/apps/odoo/bnconfig --disable_banner 1

Restart the Web server:
sudo /opt/odoo/ctlscript.sh restart apache